# config/ingest.conf
# Upstart Job to execute the Spring Boot application on startup
description "ingest"
start on startup
task
exec /vagrant/ingest/config/spring-start.sh
